# services/gunicorn.service
[Unit]
Description=lookaway gunicorn daemon
After=network.target

[Service]
User=lookaway
Group=lookaway
WorkingDirectory=/home/lookaway/lookaway-env/lookaway
ExecStart=/home/lookaway/lookaway_start

[Install]
WantedBy=multi-user.target
